# systemd/notify-prayer.timer
[Unit]
Description=Prayer notification every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=1s

[Install]
WantedBy=timers.target
