# systemd/notify-prayer.service
[Unit]
Description=Prayer Time Notification

[Service]
Type=oneshot
ExecStart=%h/.config/niri/scripts/notify-prayer
